# bin/gentoo-initscript
#!/sbin/runscript

depend() {
	need net
	use logger
}

start() {
	ebegin "Starting ${SVCNAME}"

	god --no-syslog
	
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"

	god terminate

	eend $?
}

